# ndside/msys64/var/lib/pacman/local/switch-mesa-18.2_beta2-1/files
%FILES%
opt/
opt/devkitpro/
opt/devkitpro/portlibs/
opt/devkitpro/portlibs/switch/
opt/devkitpro/portlibs/switch/include/
opt/devkitpro/portlibs/switch/include/EGL/
opt/devkitpro/portlibs/switch/include/EGL/egl.h
opt/devkitpro/portlibs/switch/include/EGL/eglext.h
opt/devkitpro/portlibs/switch/include/EGL/eglextchromium.h
opt/devkitpro/portlibs/switch/include/EGL/eglmesaext.h
opt/devkitpro/portlibs/switch/include/EGL/eglplatform.h
opt/devkitpro/portlibs/switch/include/GL/
opt/devkitpro/portlibs/switch/include/GL/gl.h
opt/devkitpro/portlibs/switch/include/GL/gl_mangle.h
opt/devkitpro/portlibs/switch/include/GL/glcorearb.h
opt/devkitpro/portlibs/switch/include/GL/glext.h
opt/devkitpro/portlibs/switch/include/GLES/
opt/devkitpro/portlibs/switch/include/GLES/egl.h
opt/devkitpro/portlibs/switch/include/GLES/gl.h
opt/devkitpro/portlibs/switch/include/GLES/glext.h
opt/devkitpro/portlibs/switch/include/GLES/glplatform.h
opt/devkitpro/portlibs/switch/include/GLES2/
opt/devkitpro/portlibs/switch/include/GLES2/gl2.h
opt/devkitpro/portlibs/switch/include/GLES2/gl2ext.h
opt/devkitpro/portlibs/switch/include/GLES2/gl2platform.h
opt/devkitpro/portlibs/switch/include/GLES3/
opt/devkitpro/portlibs/switch/include/GLES3/gl3.h
opt/devkitpro/portlibs/switch/include/GLES3/gl31.h
opt/devkitpro/portlibs/switch/include/GLES3/gl32.h
opt/devkitpro/portlibs/switch/include/GLES3/gl3ext.h
opt/devkitpro/portlibs/switch/include/GLES3/gl3platform.h
opt/devkitpro/portlibs/switch/include/KHR/
opt/devkitpro/portlibs/switch/include/KHR/khrplatform.h
opt/devkitpro/portlibs/switch/lib/
opt/devkitpro/portlibs/switch/lib/libEGL.a
opt/devkitpro/portlibs/switch/lib/libEGL.la
opt/devkitpro/portlibs/switch/lib/libglapi.a
opt/devkitpro/portlibs/switch/lib/libglapi.la
opt/devkitpro/portlibs/switch/lib/libGLESv1_CM.a
opt/devkitpro/portlibs/switch/lib/libGLESv1_CM.la
opt/devkitpro/portlibs/switch/lib/libGLESv2.a
opt/devkitpro/portlibs/switch/lib/libGLESv2.la
opt/devkitpro/portlibs/switch/lib/pkgconfig/
opt/devkitpro/portlibs/switch/lib/pkgconfig/egl.pc
opt/devkitpro/portlibs/switch/lib/pkgconfig/glesv1_cm.pc
opt/devkitpro/portlibs/switch/lib/pkgconfig/glesv2.pc
opt/devkitpro/portlibs/switch/share/
opt/devkitpro/portlibs/switch/share/drirc.d/
opt/devkitpro/portlibs/switch/share/drirc.d/00-mesa-defaults.conf

